# service/gs_maintenance/run
#!/bin/bash
source /opt/gemstone/product/seaside/defSeaside
$GEMSTONE/bin/waitstone $GEMSTONE_NAME 2
exec setuidgid DAEMONTOOLS_USER $GEMSTONE/seaside/bin/maintenanceGemSea30_oodb31
